# contrib/systemd/pushbroom.service
[Unit]
Description=Remove old files

[Service]
Type=oneshot
User=%I
ExecStart=/usr/local/bin/pushbroom
