# contrib/systemd/pushbroom.service
[Unit]
Description=Pushbroom

[Service]
User=%I
Restart=always
RestartSec=3600s
ExecStart=/usr/local/bin/pushbroom
