# contrib/systemd/pushbroom.timer
[Unit]
Description=Remove old files

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
