# crontab
# Crontab for Janitor
0  */6  *  *  * $HOME/.janitor.sh
